# supervisor.conf
[program:YuBlog]
# 启动命令入口
command=/project/YuBlog/venv/bin/uwsgi /project/YuBlog/config.ini

# 命令程序所在目录
directory=/project/YuBlog
# 运行命令的用户名
user=root

autostart=true
autorestart=true
# 日志地址
stdout_logfile=/project/YuBlog/logs/uwsgi_supervisor.log